# systemd-service/bot_lastfm_pg_monthly.timer
[Unit]
Description=bot_lastfm_pg_monthly.timer

[Timer]
OnCalendar=*-*-1 09:30:00
Persistent=true

[Install]
WantedBy=timers.target
