# systemd-service/lastfm_pg_weekly.timer
[Unit]
Description=lastfm_pg_weekly.timer

[Timer]
OnCalendar=*-*-* 10:00:00
Persistent=true

[Install]
WantedBy=timers.target
