# systemd-service/lastfm_pg_semiannually.timer
[Unit]
Description=lastfm_pg_semiannually.timer

[Timer]
OnCalendar=semiannually
Persistent=true

[Install]
WantedBy=timers.target
