# systemd-service/lastfm_pg_monthly.timer
[Unit]
Description=lastfm_pg_monthly.timer

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
