# systemd-service/twitter_bot_lastfm_pg_yearly.timer
[Unit]
Description=twitter_bot_lastfm_pg_yearly.timer

[Timer]
OnCalendar=*-1-1 08:30:00
Persistent=true

[Install]
WantedBy=timers.target
